# config/redis_ilia.service
[Unit]
Description=Redis
After=syslog.target

[Service]
ExecStart=/usr/local/bin/redis-server /simurgh/config/redis.conf
RuntimeDirectory=redis
RestartSec=5s
Restart=on-success
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
