# config/simurgh.backup.timer
[Unit]
Description=Simurgh Backup Timer

[Timer]
OnUnitActiveSec=3d
OnBootSec=1m

[Install]
WantedBy=timers.target
