# config/simurgh.backup.timer
[Unit]
Description=Simurgh Backup Timer

[Timer]
OnUnitActiveSec=1h
OnBootSec=1m

[Install]
WantedBy=timers.target
