# config/simurgh.backup.timer
[Unit]
Description=Simurgh Backup Timer

[Timer]
OnUnitActiveSec=10m
OnBootSec=10m

[Install]
WantedBy=timers.target
